# etc/services.d/mastodon-sidekiq/run
#!/usr/bin/with-contenv bash
su mastodon -c "cd /opt/mastodon/live/ && bundle exec sidekiq -c 5 -q default -q mailers -q pull -q push"