# etc/services.d/mastodon-web/run
#!/usr/bin/with-contenv bash
su mastodon -c "cd /opt/mastodon/live/ && bundle exec puma -C config/puma.rb"