# etc/services.d/mastodon-streaming/run
#!/usr/bin/with-contenv bash
su mastodon -c "cd /opt/mastodon/live/ && yarn run start"